# polly.service
[Unit]
Description=Polly Discord Poll Bot
After=network.target

[Service]
Type=exec
User=xyn0th
Group=xyn0th
WorkingDirectory=/home/xyn0th/polly
Environment=PATH=/home/xyn0th/polly/.venv/bin
EnvironmentFile=/home/xyn0th/polly/.env
ExecStart=/usr/bin/uv run python -m polly.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
